# hsg-canvas.service
[Unit]
Description=HSG Canvas SRS Server
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=hsg
Group=hsg
WorkingDirectory=/home/hsg/srs_server
ExecStart=/home/hsg/srs_server/start.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=DISPLAY=:0
Environment=HOME=/home/hsg

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/hsg/srs_server /tmp
ProtectHome=true

[Install]
WantedBy=multi-user.target